# Ansible/roles/install_brc4/templates/brc4.service.j2
[Unit]
Description=BRC4 Service

[Service]
Type=simple
WorkingDirectory=/opt/brc4
ExecStart=/bin/bash /opt/brc4/brc4.sh
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
