# Ansible/roles/install_cobaltstrike/templates/cs.service.j2
[Unit]
Description=Cobalt Strike Service

[Service]
Type=simple
WorkingDirectory=/opt/cs/cobaltstrike
ExecStart=/bin/bash cs.sh
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
