# Ansible/roles/install_havoc/templates/havoc.service.j2
[Unit]
Description=Havoc Service

[Service]
Type=simple
WorkingDirectory=/opt/havoc
ExecStart=/bin/bash havoc.sh
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
